# docker/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:rosboard]
command=/opt/supervisord/launch_rosboard.sh

[program:rosbridge]
command=/opt/supervisord/launch_rosbridge.sh

[program:microxrceagent]
command=MicroXRCEAgent udp4 -p 8888

[program:px4]
command=/opt/supervisord/launch_px4.sh

[program:control_velocity]
command=/opt/supervisord/launch_velocity.sh
